# logs.service
[Unit]
Description=Start logrotate

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate -v /etc/logrotate.conf -s /data/status